# deploy/aimusicmed-backup-requests.path
[Unit]
Description=Watch AIMusicMed backup request queue

[Path]
PathExistsGlob=/opt/aimusicmed/backups/requests/pending/*.json
Unit=aimusicmed-backup-requests.service

[Install]
WantedBy=multi-user.target
